# contrib/singularity.service
[Unit]
Description=Run Singularity and reload PDNS Lua script

[Service]
Type=oneshot
ExecStart=/opt/singularity/update.sh
