# singularity.timer
[Unit]
Description=Run Singularity daily

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target
